# voice-keyboard.service
[Unit]
Description=Voice Keyboard - Speech to Text Input
Documentation=https://github.com/deepgram/voice-keyboard-linux
After=graphical-session.target sound.target

[Service]
Type=simple
ExecStart=/bin/bash -c 'XDG_RUNTIME_DIR=${XDG_RUNTIME_DIR} exec /home/andycinquin/.local/bin/voice-keyboard-wrapper'
Restart=always
RestartSec=5s

# Environment variables
Environment="XDG_RUNTIME_DIR=/run/user/1000"
Environment="PULSE_SERVER=unix:/run/user/1000/pulse/native"
Environment="RUST_LOG=info"

# Logs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=voice-keyboard

[Install]
WantedBy=default.target
